# scripts/expo-metro.service
# expo-metro.service — systemd unit for Expo Metro bundler
#
# Instalación en el servidor de producción:
#   sudo cp scripts/expo-metro.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable expo-metro
#   sudo systemctl start expo-metro
#
# Logs:
#   sudo journalctl -u expo-metro -f
#
# Estado:
#   sudo systemctl status expo-metro

[Unit]
Description=Expo Metro Bundler — Web + Expo Go
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=root
WorkingDirectory=/root/react-native-strength-training

Environment=EXPO_PACKAGER_PROXY_URL=https://entrenamentua.musikak.com
Environment=REACT_NATIVE_PACKAGER_HOSTNAME=entrenamentua.musikak.com

# Puerto 8082 para no competir con expo-entrenamentua.service (8081)
ExecStart=/root/react-native-strength-training/node_modules/.bin/expo start --web --port 8082 --clear

# Timeout generoso para startup con cache limpia
TimeoutStartSec=120

# Timeout para evitar zombie processes si el stop se cuelga
TimeoutStopSec=30

# Mata todo el control group al hacer stop (evita huérfanos)
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes

# Standard output goes to journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
