# scripts/expo-metro-web.service
# expo-metro-web.service — systemd unit for Expo Metro bundler (web)
#
# Servidor web (port 8083). Corre junto con expo-metro.service (port 8082, Expo Go).
#
# Instalación:
#   sudo cp scripts/expo-metro-web.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable expo-metro-web
#   sudo systemctl start expo-metro-web
#
# Logs:
#   sudo journalctl -u expo-metro-web -f

[Unit]
Description=Expo Metro Bundler — Web (entrenamentua.musikak.com)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/react-native-strength-training

Environment=EXPO_PACKAGER_PROXY_URL=https://entrenamentua.musikak.com
Environment=REACT_NATIVE_PACKAGER_HOSTNAME=entrenamentua.musikak.com

ExecStart=/usr/bin/npx expo start --web --port 8083

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
